# system/systemd/companion-ssh-apply.service
[Unit]
Description=Apply Companion-requested SSH access changes
Documentation=https://github.com/jklingberg/partybox-companion

# Never started at boot -- only ever triggered on demand by companion.service
# (running as the unprivileged 'companion' user) over D-Bus, authorized by a
# polkit rule scoped to exactly this unit name. See ADR-042 and
# image/config/companion-ssh-apply.sh.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/companion-ssh-apply
