# system/systemd/companion.service
[Unit]
Description=PartyBox Companion
Documentation=https://github.com/jklingberg/partybox-companion
After=network-online.target bluetooth.service NetworkManager.service
Wants=network-online.target bluetooth.service NetworkManager.service

[Service]
Type=simple
User=companion
Group=bluetooth

# Production path for persistent data — always set, cannot be overridden accidentally.
# The Portal writes config.json here; systemd creates the directory before ExecStart.
Environment=COMPANION_DATA_DIR=/var/lib/companion

# Operator overrides: copy system/systemd/companion.env to /etc/companion/companion.env
# and uncomment the settings you want to change. The file is silently ignored if absent.
EnvironmentFile=-/etc/companion/companion.env

# Reset the HCI controller before starting to clear wedged state from prior failed
# connections. The '+' prefix runs this as root regardless of User=companion.
ExecStartPre=+/usr/bin/hciconfig hci0 reset

ExecStart=/opt/partybox-companion/bin/partybox-companion

Restart=on-failure
RestartSec=5
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=companion

# StateDirectory creates /var/lib/companion owned by companion before ExecStart.
StateDirectory=companion

# RuntimeDirectory creates /run/companion (tmpfs-backed, owned by companion) and
# sets XDG_RUNTIME_DIR=/run/companion. This is the correct location for PipeWire
# sockets, WirePlumber state, and any future user-session runtime artefacts.
# Having this in place means enabling PipeWire for the companion user requires no
# change to this unit — see ADR-019 companion user section.
RuntimeDirectory=companion
RuntimeDirectoryMode=0700

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Port 80 binding: CAP_NET_BIND_SERVICE allows the companion user to bind port 80
# without running as root. AmbientCapabilities grants it to the process at exec;
# CapabilityBoundingSet restricts the process from acquiring any other capability.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
