# system/systemd/companion.service
[Unit]
Description=PartyBox Companion
Documentation=https://github.com/jklingberg/partybox-companion
After=network-online.target bluetooth.service
Wants=network-online.target bluetooth.service

[Service]
Type=simple
User=companion
Group=bluetooth

# Production path for persistent data — always set, cannot be overridden accidentally.
# The Portal writes config.json here; systemd creates the directory before ExecStart.
Environment=COMPANION_DATA_DIR=/var/lib/companion

# Operator overrides: copy system/systemd/companion.env to /etc/companion/companion.env
# and uncomment the settings you want to change. The file is silently ignored if absent.
EnvironmentFile=-/etc/companion/companion.env

ExecStart=/opt/partybox-companion/bin/partybox-companion

Restart=on-failure
RestartSec=5
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=companion

# StateDirectory creates /var/lib/companion owned by companion before ExecStart.
StateDirectory=companion

# RuntimeDirectory creates /run/companion (tmpfs-backed, owned by companion) and
# sets XDG_RUNTIME_DIR=/run/companion. This is the correct location for PipeWire
# sockets, WirePlumber state, and any future user-session runtime artefacts.
# Having this in place means enabling PipeWire for the companion user requires no
# change to this unit — see ADR-019 companion user section.
RuntimeDirectory=companion
RuntimeDirectoryMode=0700

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Port 80 binding: the appliance exposes the Portal on port 80 in production.
# The mechanism (AmbientCapabilities, reverse proxy, or iptables) is determined
# at image build time (M14). Set COMPANION_PORT=80 in companion.env once decided.

[Install]
WantedBy=multi-user.target
