# system/systemd/partyboxd.service
[Unit]
Description=PartyBox Companion Daemon
Documentation=https://github.com/jklingberg/partybox-companion
After=bluetooth.target network-online.target avahi-daemon.service
Wants=bluetooth.target network-online.target
Requires=dbus.socket avahi-daemon.service

[Service]
Type=simple
User=partybox
Group=partybox
SupplementaryGroups=bluetooth audio

ExecStart=/usr/local/bin/partyboxd --config /etc/partyboxd/partyboxd.toml
Restart=always
RestartSec=10

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/partyboxd /run/partyboxd

RuntimeDirectory=partyboxd
StateDirectory=partyboxd
LogsDirectory=partyboxd

[Install]
WantedBy=multi-user.target
